# internal/assets/templates/projects/generic-service/deploy/systemd/app.service.tmpl
[Unit]
Description={{.App}} service
After=network.target

[Service]
User={{.ServiceUser}}
Group={{.ServiceGroup}}
WorkingDirectory=/opt/{{.App}}
# EnvironmentFile loads KEY=VALUE pairs only. For YAML/INI/properties config,
# pass the path via an env var or argument in ExecStart instead.
EnvironmentFile=-/etc/{{.App}}/app.env
# Update ExecStart to point to your binary or interpreter:
#   Compiled binary:  ExecStart=/opt/{{.App}}/{{.App}}
#   Python:           ExecStart=/usr/bin/python3 /opt/{{.App}}/main.py
#   Shell script:     ExecStart=/opt/{{.App}}/start.sh
ExecStart=/opt/{{.App}}/{{.App}}
Restart=always
RestartSec=10
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/{{.App}}

[Install]
WantedBy=multi-user.target
