# internal/assets/templates/projects/node-service/deploy/systemd/app.service.tmpl
[Unit]
Description={{.App}} service
After=network.target

[Service]
User={{.ServiceUser}}
Group={{.ServiceGroup}}
WorkingDirectory=/opt/{{.App}}
# EnvironmentFile loads KEY=VALUE pairs only. For YAML/INI/properties config,
# pass the path via an env var or argument in ExecStart instead.
EnvironmentFile=-/etc/{{.App}}/app.env
# Nuxt 3: /usr/bin/node /opt/{{.App}}/server/index.mjs
# Express/custom: /usr/bin/node /opt/{{.App}}/index.js
ExecStart=/usr/bin/node /opt/{{.App}}/server/index.mjs
Restart=always
RestartSec=10
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/{{.App}}

[Install]
WantedBy=multi-user.target
