# internal/assets/templates/projects/multi-service/deploy/systemd/backend.service.tmpl
[Unit]
Description={{.App}} backend
After=network.target

[Service]
User={{.ServiceUser}}
Group={{.ServiceGroup}}
WorkingDirectory=/opt/{{.App}}/backend
# EnvironmentFile loads KEY=VALUE pairs only. For YAML/INI/properties config,
# pass the path via an env var or argument in ExecStart instead.
EnvironmentFile=-/etc/{{.App}}/app.env
ExecStart=/usr/bin/java $JAVA_OPTS -jar /opt/{{.App}}/backend/backend.jar
Restart=always
RestartSec=10
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/{{.App}}

[Install]
WantedBy=multi-user.target
